# pickup/pickup.crontab
# Poll once every 5 minutes
*/5 * * * * /projects/cicd/deaddropdelivery/pickup/pickup.sh crazychenz/vinnie.work
